# ic-os/components/hostos-scripts/monitoring/custom-metrics.service
[Unit]
Description=Custom HostOS metrics
After=node_exporter.service

[Service]
Type=oneshot
ExecStart=/opt/ic/bin/custom-metrics.sh

[Install]
WantedBy=multi-user.target
