# ic-os/components/monitoring/hostos/monitor-nvme.timer
[Unit]
Description=Monitor NVMe SMART log info on a regular basis

[Timer]
Unit=monitor-nvme.service
OnCalendar=hourly

[Install]
WantedBy=timers.target
