# ic-os/components/monitoring/hostos/setup-node-gen-status.service
[Unit]
Description=Write node generation status
Before=node_exporter.service

[Service]
Type=oneshot
ExecStart=/opt/ic/bin/hostos_tool set-hardware-gen-metric
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
