# ic-os/boundary-guestos/rootfs/etc/systemd/system/vector.service.d/override.conf
[Service]
Environment="ELASTICSEARCH_URL=https://elasticsearch.mercury.dfinity.systems"
