# ic-os/components/monitoring/setup-node-gen-status.service
[Unit]
Description=Write node generation status
After=node_exporter.service

[Service]
Type=oneshot
ExecStart=/opt/ic/bin/guestos_tool set-hardware-gen-metric
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
