# ic-os/components/monitoring/node_exporter/node_exporter.service
[Unit]
Description=Node Exporter
After=syslog.target network.target

[Service]
User=node_exporter
Group=node_exporter
EnvironmentFile=/etc/default/node_exporter
ExecStart=/usr/local/bin/node_exporter $OPTIONS
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
