# ic-os/rootfs/boundary-guestos/etc/systemd/system/bootstrap-ic-node.service
[Unit]
Description=Bootstrap the IC node
Requires=var-log.mount
After=var-log.mount

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/opt/ic/bin/bootstrap-ic-node.sh
