# ic-os/components/monitoring/guestos/boot-metrics/boot-metrics.timer
[Unit]
Description=Run boot-metrics service 10 minutes after boot

[Timer]
OnBootSec=10min

[Install]
WantedBy=timers.target
