# ic-os/components/monitoring/guestos/boot-metrics/boot-metrics.service
[Unit]
Description=Collect boot metrics for Prometheus

[Service]
Type=oneshot
ExecStart=/opt/ic/bin/boot-metrics.sh
