# ic-os/guestos/rootfs/etc/systemd/system/reload_nftables.path
[Unit]
Description=Watches the nftables configuration for changes and reloads it
Before=nftables.service

[Install]
WantedBy=multi-user.target

[Path]
PathChanged=/run/ic-node/nftables-ruleset/nftables.conf
