# ic-os/boundary-guestos/rootfs/etc/systemd/system/journalbeat-restart.timer
[Unit]
Description=Daily restart of journalbeat

[Timer]
OnCalendar=daily
AccuracySec=12h
Persistent=true

[Install]
WantedBy=timers.target
