# ic-os/boundary-guestos/rootfs/etc/systemd/system/journalbeat-restart.service
[Unit]
Description=Restart journalbeat

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart journalbeat
