# ic-os/components/monitoring/guestos/metrics_tool.timer
[Unit]
Description=Collect custom metrics every minute

[Timer]
OnBootSec=60s
OnUnitActiveSec=60s
Unit=metrics_tool.service

[Install]
WantedBy=timers.target
