# ic-os/components/boundary-guestos/etc/systemd/system/nftables.service.d/override.conf
[Unit]
After=setup-nftables.service
BindsTo=setup-nftables.service
