# ic-os/boundary-guestos/rootfs/etc/systemd/system/nginx.service.d/override.conf
[Service]
LimitNOFILE=524288
LogRateLimitIntervalSec=1ms
LogRateLimitBurst=1000
ExecStop=
ExecStop=/opt/ic/bin/nginx-hbpause-stop.sh
