# ic-os/components/boundary-guestos/etc/systemd/system/setup-crowdsec.service
[Unit]
Description=Create configuration for Crowdsec
DefaultDependencies=no
After=bootstrap-ic-node.service
Requires=bootstrap-ic-node.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/opt/ic/bin/setup-crowdsec.sh
