# ic-os/components/misc/log-config/log-config.service
[Unit]
Description=Log config partition
After=update-hostos-config.service
Wants=update-hostos-config.service

[Service]
Type=oneshot
ExecStart=/opt/ic/bin/log-config.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target