# ic-os/components/misc/guestos-recovery/guestos-recovery-upgrader/guestos-recovery-upgrader.service
[Unit]
Description=Recovery Upgrader
ConditionKernelCommandLine=recovery=1
After=network-online.target
Requires=network-online.target
Before=guestos.service

[Service]
Type=oneshot
ExecStart=/opt/ic/bin/guestos-recovery-upgrader.sh
StandardOutput=journal+console
StandardError=journal+console
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
