# ic-os/boundary-guestos/rootfs/etc/cron.d/journalbeat-restart
0 * * * * root systemctl restart journalbeat
