# ic-os/components/ic/ic-btc-adapter/ic-doge-mainnet-adapter.service
[Unit]
Description=IC Dogecoin Mainnet Adapter Service Provider
After=bootstrap-ic-node.service
Wants=bootstrap-ic-node.service
After=network-online.target
Wants=network-online.target
Requires=ic-doge-mainnet-adapter.socket
StartLimitIntervalSec=0

[Service]
User=ic-replica
ExecStartPre=+/opt/ic/bin/generate-doge-adapter-config.sh -m -o /run/ic-node/config/ic-doge-mainnet-adapter.json5
ExecStart=/opt/ic/bin/ic-btc-adapter /run/ic-node/config/ic-doge-mainnet-adapter.json5
Restart=always

[Install]
WantedBy=ic-replica.service
Also=ic-doge-mainnet-adapter.socket
