# ic-os/components/init/setup-lvs/setup-lvs.service
[Unit]
Description=Set up logical volumes on first boot
DefaultDependencies=no
Requires=open-store.service
After=open-store.service

[Install]
WantedBy=local-fs.target

[Service]
StandardOutput=journal+console
StandardError=journal+console
Type=oneshot
RemainAfterExit=true
ExecStart=/opt/ic/bin/setup-lvs.sh
