# ic-os/components/networking/nftables/reload_nftables.service
[Unit]
Description=Reload nftables when the configuration changes

[Service]
Type=oneshot
ExecStartPre=/usr/sbin/nft flush ruleset
ExecStart=systemctl reload nftables.service
