# ic-os/guestos/containerfs/etc/systemd/system/systemd-logind.service.d/override.conf
# This removes delay incurred at ssh login when running as
# unprivileged docker containers (see bug
# https://github.com/systemd/systemd/issues/15408)
[Service]
ProtectHostname=no
