# ic-os/components/misc/vsock/vsock-agent.service
[Unit]
Description=VSOCK agent daemon

[Service]
ExecStart=/opt/ic/bin/vsock_host
Restart=on-failure

[Install]
WantedBy=multi-user.target
