# ic-os/components/hostos-scripts/generate-guestos-config/generate-guestos-config.service
[Unit]
Description=Generate GuestOS configuration
Before=systemd-networkd.service
# TODO(NODE-1618):
After=update-config.service
Wants=update-config.service
RequiresMountsFor=/var

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/opt/ic/bin/config generate-guest-vm-config
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
