# ic-os/components/monitoring/guestos/ipv4-connectivity-check/ipv4-connectivity-check.service
[Unit]
Description=Check IPv4 connectivity

[Install]
WantedBy=multi-user.target

[Service]
Type=oneshot
ExecStart=/opt/ic/bin/ipv4-connectivity-check.sh

# Disable systemd start and stop logs
LogLevelMax=1
SyslogLevel=2
