# ic-os/components/boundary-guestos/etc/systemd/system/nginx.service.d/override.conf
[Unit]
After=setup-nginx.service
BindsTo=setup-nginx.service

[Service]
LimitNOFILE=524288
LogRateLimitIntervalSec=1ms
LogRateLimitBurst=1000
